# update-ip-cf.timer
[Unit]
Description=Run update-ip-cf hourly

[Install]
WantedBy=timers.target

[Timer]
OnCalendar=hourly
Persistent=true
